# iso/worker/airootfs/etc/systemd/system/docker.service.d/10-partitions.conf
[Unit]
BindsTo=make-partitions.service
After=make-partitions.service
