# iso/worker/airootfs/etc/systemd/system/chronyd.service.d/10-settings.conf
[Unit]
Requires=apply-settings.service
After=apply-settings.service

[Service]
ExecStartPre=/usr/bin/chronyd -q -u chrony
TimeoutStartSec=infinity
